# docker/FlaskFileSystem/etc/app.ini
[program:app]
command=flask run -h 0.0.0.0
autostart=true
autorestart=true
environment=FLASK_APP="file_system.py"
stdout_logfile=logs/app_out.log
stderr_logfile=logs/app_err.log
